# nas_ports/freenas/freenas-networkd/pkg-plist
lib/networkd/src/main.py
lib/networkd/src/main.pyc
etc/rc.d/networkd
sbin/networkd
