# nas_ports/freenas/freenas-schedulerd/pkg-plist
lib/schedulerd/src/main.py
lib/schedulerd/src/main.pyc
etc/rc.d/schedulerd
sbin/schedulerd
