# nas_ports/freenas/freenas-schedulerd/pkg-plist
/usr/local/lib/schedulerd/src/main.py
/usr/local/etc/rc.d/schedulerd
/usr/local/sbin/schedulerd
