# nas_ports/freenas/freenas-crashd/pkg-plist
lib/crashd/src/main.py
lib/crashd/src/main.pyc
etc/rc.d/crashd
sbin/crashd
