# nas_ports/freenas/freenas-dispatcher/pkg-plist
lib/node_modules/dispatcher-client/package.json
lib/node_modules/dispatcher-client/lib/DispatcherClient.js
lib/node_modules/dispatcher-client/lib/EntitySubscriber.js
lib/dispatcher/src/frontend/templates/nav.html
lib/dispatcher/src/frontend/templates/tasks.html
lib/dispatcher/src/frontend/templates/events.html
lib/dispatcher/src/frontend/templates/rpc.html
lib/dispatcher/src/frontend/templates/term.html
lib/dispatcher/src/frontend/templates/stats.html
lib/dispatcher/src/frontend/templates/apidoc/index.html
lib/dispatcher/src/frontend/templates/apidoc/nav.html
lib/dispatcher/src/frontend/templates/apidoc/rpc.html
lib/dispatcher/src/frontend/templates/apidoc/events.html
lib/dispatcher/src/frontend/templates/apidoc/tasks.html
lib/dispatcher/src/frontend/templates/apidoc/schemas.html
lib/dispatcher/src/frontend/static/middleware.js
lib/dispatcher/src/frontend/static/emitter.js
lib/dispatcher/src/frontend/static/term.js
lib/dispatcher/src/frontend/static/jquery.jsonview.js
lib/dispatcher/src/frontend/static/jquery.jsonview.css
lib/dispatcher/src/frontend/static/main.css
libexec/dispatcher/notify-pam
libexec/dispatcher/ssh_transport_catcher
lib/dispatcher/src/lib/sendzfs.so
etc/rc.d/dispatcher
sbin/dispatcherd
sbin/dispatcherctl
sbin/sdsinit
libexec/taskworker
etc/middleware.conf
