# nas_ports/freenas/freenas-dispatcher/pkg-plist
%%PYTHON_SITELIBDIR%%/dispatcher_client-10.0-py2.7.egg-info
%%PYTHON_SITELIBDIR%%/dispatcherclient.py
%%PYTHON_SITELIBDIR%%/dispatcherclient.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/rpc.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/rpc.py
%%PYTHON_SITELIBDIR%%/dispatcher/spawn_thread.py
%%PYTHON_SITELIBDIR%%/dispatcher/spawn_thread.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/client_transport.py
%%PYTHON_SITELIBDIR%%/dispatcher/client_transport.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/client.py
%%PYTHON_SITELIBDIR%%/dispatcher/__init__.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/__init__.py
%%PYTHON_SITELIBDIR%%/dispatcher/client.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/validator.py
%%PYTHON_SITELIBDIR%%/dispatcher/validator.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/jsonenc.py
%%PYTHON_SITELIBDIR%%/dispatcher/jsonenc.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/shell.py
%%PYTHON_SITELIBDIR%%/dispatcher/shell.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/file.py
%%PYTHON_SITELIBDIR%%/dispatcher/file.pyc
lib/dispatcher/src/frontend/templates/nav.html
lib/dispatcher/src/frontend/templates/tasks.html
lib/dispatcher/src/frontend/templates/events.html
lib/dispatcher/src/frontend/templates/rpc.html
lib/dispatcher/src/frontend/templates/term.html
lib/dispatcher/src/frontend/templates/stats.html
lib/dispatcher/src/frontend/templates/apidoc/index.html
lib/dispatcher/src/frontend/templates/apidoc/nav.html
lib/dispatcher/src/frontend/templates/apidoc/rpc.html
lib/dispatcher/src/frontend/templates/apidoc/events.html
lib/dispatcher/src/frontend/templates/apidoc/tasks.html
lib/dispatcher/src/frontend/templates/apidoc/schemas.html
lib/dispatcher/src/frontend/static/middleware.js
lib/dispatcher/src/frontend/static/emitter.js
lib/dispatcher/src/frontend/static/term.js
lib/dispatcher/src/frontend/static/jquery.jsonview.js
lib/dispatcher/src/frontend/static/jquery.jsonview.css
lib/dispatcher/src/frontend/static/main.css
libexec/dispatcher/notify-pam
libexec/dispatcher/ssh_transport_catcher
etc/rc.d/dispatcher
sbin/dispatcherd
sbin/dispatcherctl
libexec/taskproxy
etc/middleware.conf
