# src/middleware-files/usr/local/bin/ipfs
#!/bin/sh
. /etc/rc.subr
load_rc_config "ipfs_go"
[ -n "$ipfs_go_path" ] && export IPFS_PATH=$ipfs_go_path
exec /usr/local/bin/ipfs-go $*
