# src/fnstatd/etc/fnstatd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: fnstatd
# REQUIRE: datastore dispatcher
# BEFORE: collectd

. /etc/rc.subr

name="fnstatd"
rcvar="fnstatd_enable"
pidfile="/var/run/fnstatd.pid"
command="/usr/sbin/daemon"
start_precmd="prestart"

prestart()
{
    rc_flags="-P ${pidfile} -f -r /usr/local/sbin/crash-wrapper /usr/local/sbin/fnstatd ${rc_flags}"
}

load_rc_config $name
run_rc_command $*
